# xenon-slurm-15.08/service/slurmdbd/run
#!/bin/sh

# Make sure slurm controller is restarted before slurmdbd is started, to force controller to reconnect
killall slurmctld

exec /usr/sbin/slurmdbd -D