# xenon-slurm-15.08/service/slurmctld/run
#!/bin/sh

# Create mycluster tables in mysql via slurmdbd
sacctmgr -i add cluster mycluster

exec /usr/sbin/slurmctld -D -c