# xenon-ssh/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -De

[program:atd]
command=/usr/sbin/atd -f -b 5


