# xenon-slurm/xenon-slurm-ssh/supervisor/sshd.conf
[program:sshd]
directory=/usr/local/
command=/usr/sbin/sshd -De
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd.log
